# solaredge/rootfs/etc/services.d/capturer/run
#!/usr/bin/with-contenv bashio

bashio::log.notice "Welcome to Solaredge Data Capture!"

capture_path="/data/traffic/"
se_logger_service="/opt/se-logger/liveupdate.py"

#Create capture path
mkdir -p $capture_path
chmod a+x $se_logger_service

#Write the SE-Logger's config file out
bashio::log.info "Writing config file..."
cat > /opt/se-logger/config.py <<- EOF
#!/usr/bin/env python
inverter_private_key = b'$(bashio::config 'inverter_private_key')'
db_user = "$(bashio::config 'database_user')"
db_pass = "$(bashio::config 'database_password')"
db_name = "$(bashio::config 'database_name')"
db_host = "$(bashio::config 'database_host')"
db_port = $(bashio::config 'database_port')
EOF

INTERFACE="$(bashio::config 'capture_interface')"
FILTER=tcp
PREFIX=solaredge-

rm -f /data/tcpdump.log

# Remove any non-zero files. We don't need them
find -name '/data/traffic/*' -size 0 -delete

# wait for the time to get set
while [ `date -u +%Y` == "1970" ]
do
	bashio::log.info "Waiting for correct time..."
	sleep 1
done

bashio::log.info "Setting up bridge network.."
# Create interfaces
ip link add name $INTERFACE type bridge || true
ip link add link eth0 name eth0.11 type vlan id 11 || true
ip link add link eth0 name eth0.12 type vlan id 12 || true
# Bring them up
bashio::log.info "Bringin interfaces UP.."
ip link set $INTERFACE up
ip link set eth0.11 up
ip link set eth0.12 up
# Enslave them to the bridge
bashio::log.info "Enslaving bridge interfaces.."
ip link set eth0.11 master $INTERFACE
ip link set eth0.12 master $INTERFACE
# Enable traffic forwarding in the firewall
bashio::log.info "Setup traffic forwarding in netfilter.."

# old school
# iptables -D DOCKER-USER -i $INTERFACE -o $INTERFACE -j ACCEPT || true
# iptables -I DOCKER-USER -i $INTERFACE -o $INTERFACE -j ACCEPT

# new school
nft insert rule ip filter DOCKER-USER iifname $INTERFACE oifname $INTERFACE counter accept comment \"solaredgebridge\"
bashio::log.info "done"

# sleep infinity

# Start dumping the traffic on the bridge
bashio::log.info "Starting Solaredge Bridge.."
/usr/bin/stdbuf -i0 -o0 -e0 tcpdump -i $INTERFACE -U --micro -w - ${FILTER} 2>> /data/tcpdump.log | \
	tee ${capture_path}$PREFIX`date -u +%Y%m%d%H%M%S`.pcap | \
	/usr/bin/python3 -u ${se_logger_service} - &
solaredge_pid=$!
bashio::log.info "done"

# And let it roll
bashio::log.notice "Solaredge Bridge is running!"
wait "${solaredge_pid}"
